# .circleci/root/etc/services.d/app/run
#!/usr/bin/with-contenv sh

cd /app

exec node app.js;
